# deploy/systemd/bnc-ptpmon.service
[Unit]
Description=BKG NTRIP Client (BNC) — float PPP cross-engine validator
Documentation=https://igs.bkg.bund.de/ntrip/bnc
After=ntrip-caster-local-ptpmon.service network-online.target
Wants=ntrip-caster-local-ptpmon.service network-online.target

[Service]
Type=simple
User=bob
Group=bob
ExecStart=/home/bob/peppar-bnc-glue/deploy/run-bnc.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# BNC + future PPP-Wizard library may load orbits into RAM.  Cap
# generously to keep the box healthy.
MemoryMax=2G

[Install]
WantedBy=multi-user.target
