# deploy/systemd/f9t-tcp-bridge-ptpmon.service
[Unit]
Description=F9T (kernel GNSS char device) → local TCP RTCM 3 bridge
Documentation=https://github.com/bobvan/peppar-bnc-glue
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# /dev/gnss0 is root-only (crw------- root:root) on E810.  Run as
# root for the smoke-test phase.  TODO: install a udev rule that
# makes /dev/gnss* group-readable by `dialout` or a dedicated
# `gnss` group, then drop User= to bob.
User=root
Group=root

# Bridge raw bytes (RTCM 3 stream) from the F9T's kernel GNSS char
# device to a local TCP listener.  BNC and any tcp client read
# from port 2101.
#
# Reader/broadcast architecture: ONE process reads /dev/gnss0 and
# fans the bytes out to every connected client.  This replaces the
# earlier `socat -u FILE:/dev/gnss0 TCP-LISTEN:2101,fork` pattern,
# which had a subtle bug — `fork` spawned a separate /dev/gnss0
# reader per accepted client, and multiple readers of a Linux char
# device SPLIT the byte stream rather than each receiving a copy.
# Measured 2026-04-27 on ptpmon: 2-consumer setup got 194 B/s each
# vs single-consumer 633 B/s, with bytes lost in races.
#
# Historical note on RTKLIB str2str: also doesn't work for
# /dev/gnss*.  `-in serial://` calls termios.tcgetattr() (EINVAL on
# kernel char devices), `-in file://` opens but never reads.
ExecStart=/usr/bin/python3 /home/bob/peppar-bnc-glue/deploy/f9t_broadcaster.py \
    --device /dev/gnss0 \
    --listen-host 127.0.0.1 \
    --listen-port 2101

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
