# deploy/systemd/str2str-ptpmon.service
[Unit]
Description=str2str — F9T (kernel GNSS) → local NTRIP mount
Documentation=https://github.com/PrideLab/PRIDE-PPPAR
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bob
Group=bob

# Read F9T from kernel GNSS char device (E810 AQ-mediated I2C),
# republish as RTCM 3 over a local NTRIP mount on port 2101.
# F9T must be pre-configured to emit 1077/1097/1127/1005/1230 on the
# I2C port — see config/ptpmon-rtcm.toml + deploy/configure-f9t-rtcm.sh.
#
# str2str command-line lives in config/str2str.args (gitignored,
# operator-edited from str2str.args.example).
ExecStart=/usr/bin/str2str @/home/bob/peppar-bnc-glue/config/str2str.args

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
